# pipeline/com.odcus.knowhow-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.odcus.knowhow-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <!-- Update this path to match where you copied watcher.sh -->
    <string>/Users/YOUR_USERNAME/.local/bin/knowhow-watcher.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- Update these paths to match your username -->
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/.local/logs/knowhow-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/.local/logs/knowhow-watcher.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <!-- Update this to your home directory -->
    <string>/Users/YOUR_USERNAME</string>
  </dict>
</dict>
</plist>
